# SL-CP/etc/crontabs/phred
SHELL=/bin/bash
MAILTO=phred

* * * * *  /usr/bin/perl  /home/phred/sl/SL-CP/bin/sl_cp_watch
