# SL-Linksys/src/OpenWrt-ImageBuilder-Linux-i686/build_mipsel/linux-2.4-brcm/root/usr/lib/ipkg/info/base-files.list
/
/bin/
/bin/ipcalc.sh
/bin/login
/bin/firstboot
/etc/
/etc/init.d/
/etc/init.d/S40network
/etc/init.d/S50telnet
/etc/init.d/S10boot
/etc/init.d/rcS
/etc/init.d/S60cron
/etc/init.d/S99done
/etc/init.d/S05nvram
/etc/init.d/S50httpd
/etc/hotplug.d/
/etc/hotplug.d/net/
/etc/hotplug.d/net/01-wds
/etc/hotplug.d/usb/
/etc/hotplug.d/usb/01-mount
/etc/config/
/etc/functions.sh
/etc/protocols
/etc/shells
/etc/inittab
/etc/rc.common
/etc/group
/etc/profile
/etc/hosts
/etc/passwd
/etc/openwrt_version
/etc/sysctl.conf
/etc/banner
/etc/preinit
/etc/ipkg.conf
/etc/crontabs/
/etc/resolv.conf
/rom/
/rom/note
/sbin/
/sbin/hotplug
/sbin/ifup
/sbin/mount_root
/sbin/ifdown
/sbin/hotplug.failsafe
/usr/
/usr/share/
/usr/share/udhcpc/
/usr/share/udhcpc/default.script
/usr/lib/
/usr/lib/common.awk
/usr/lib/parse-config.awk
/usr/bin/
/dev/
/jffs/
/lib/
/mnt/
/proc/
/tmp/
/www/
/var
