# SL-Linksys/src/OpenWrt-ImageBuilder-Linux-i686/build_mipsel/linux-2.4-brcm/root/usr/lib/ipkg/info/dnsmasq.list
/
/usr/
/usr/sbin/
/usr/sbin/dnsmasq
/etc/
/etc/init.d/
/etc/init.d/S60dnsmasq
/etc/dnsmasq.conf
