# SL-Linksys/src/OpenWrt-ImageBuilder-Linux-i686/build_mipsel/linux-2.4-brcm/root/usr/lib/ipkg/info/dropbear.list
/
/usr/
/usr/bin/
/usr/bin/scp
/usr/bin/ssh
/usr/bin/dbclient
/usr/bin/dropbearkey
/usr/sbin/
/usr/sbin/dropbear
/etc/
/etc/init.d/
/etc/init.d/S50dropbear
