# SL-Linksys/src/OpenWrt-ImageBuilder-Linux-i686/build_mipsel/linux-2.4-brcm/root/usr/lib/ipkg/info/iptables.list
/
/etc/
/etc/config/
/etc/config/firewall
/etc/init.d/
/etc/init.d/S35firewall
/etc/firewall.user
/usr/
/usr/lib/
/usr/lib/iptables/
/usr/lib/iptables/libipt_standard.so
/usr/lib/iptables/libipt_icmp.so
/usr/lib/iptables/libipt_tcp.so
/usr/lib/iptables/libipt_udp.so
/usr/lib/iptables/libipt_state.so
/usr/lib/iptables/libipt_limit.so
/usr/lib/iptables/libipt_multiport.so
/usr/lib/iptables/libipt_REJECT.so
/usr/lib/iptables/libipt_mark.so
/usr/lib/iptables/libipt_MARK.so
/usr/lib/iptables/libipt_TCPMSS.so
/usr/lib/iptables/libipt_SNAT.so
/usr/lib/iptables/libipt_DNAT.so
/usr/lib/iptables/libipt_MASQUERADE.so
/usr/lib/firewall.awk
/usr/sbin/
/usr/sbin/iptables
