# SL-ROBIN/package/silverlining/files/etc/init.d/sln
#!/bin/sh /etc/rc.common

# Copyright 2009 Silver Lining Networks Inc.

# /init.d/sln

START=98
ME="sln"

start() {
	logger -s -t "$ME" "starting..."

	# turn nodogsplash thread limiting
	$(uci set cp_switch.handler_1.max_threads=16)

	# force reboot on oops
	$(echo 1 > /proc/sys/kernel/panic_on_oops)

	# set the cron mod to a random time
	$(/usr/bin/microperl /usr/bin/sl_checkin)

	# turn on dns caching
	if [ -z $(grep cache /tmp/dnsmasq.conf) ] ; then
		$(echo "cache-size=4096" >> /tmp/dnsmasq.conf)

		# restart dnsmasq
		kill $(pgrep dnsmasq)
		sleep 1
		dnsmasq -C $(cat /tmp/dnsmasq.args)
	fi
}

stop() {

	$(/usr/bin/microperl /usr/bin/sl_fw_ha allstop)
	$(/usr/bin/microperl /usr/bin/sl_checkin stop)
}

case $1 in

	start)
		start
		;;

	stop)
		stop
		;;

esac

