# SL-CP/conf/crontab/sln
* * * * * /usr/bin/perl /usr/bin/sl_tunnel
