# SL-ROBIN/package/silverlining/files/etc/crontabs/nobody
*/2 * * * * /usr/bin/microperl /usr/bin/sl_fw_ha
